# deploy/app.service.example
[Unit]
Description=TraceDownloader
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Copy this file to /etc/systemd/system/tracedownloader.service and
# replace the paths below with wherever you installed the app, then:
#   systemctl daemon-reload && systemctl enable --now tracedownloader
WorkingDirectory=/opt/tracedownloader
Environment=APP_HOME=/opt/tracedownloader/data
Environment=APP_BIN_DIR=/opt/tracedownloader
Environment=APP_DEFAULT_OUTPUT=/opt/tracedownloader/download
ExecStart=/opt/tracedownloader/venv/bin/uvicorn server:app --host 127.0.0.1 --port 8686 --timeout-graceful-shutdown 5
Restart=on-failure
RestartSec=3
# Keeps an open SSE connection (/api/events) from blocking a restart.
TimeoutStopSec=15

[Install]
WantedBy=multi-user.target
